# spark-docker/slave.conf
[supervisord]
nodaemon=true

[program:spark-slave]
command=/opt/spark/sbin/start-slave.sh spark://master:7077
